# config/systemd/darkpool-collector-beat.service
[Unit]
Description=Dark Pool Collector Celery Beat
After=network.target redis.service

[Service]
Type=simple
User=avxz
Group=avxz
WorkingDirectory=/opt/darkpool_collector
Environment="PATH=/opt/darkpool_collector/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="ENV_FILE=/opt/darkpool_collector/.env.prod"
EnvironmentFile=/opt/darkpool_collector/.env.prod
ExecStart=/opt/darkpool_collector/venv/bin/celery -A config.celery.celery_app beat --loglevel=info --schedule=/var/lib/celery/beat/darkpool-beat-schedule.db --broker-connection-retry-on-startup=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 