# config/systemd/news-collector-worker.service
[Unit]
Description=News Collector Celery Worker
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/vvv-invest
Environment=PYTHONPATH=/home/ubuntu/vvv-invest
Environment=ENV_FILE=/home/ubuntu/vvv-invest/.env
ExecStart=/home/ubuntu/vvv-invest/venv/bin/celery -A config.celery.celery_app worker -Q news_queue -l info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 